# deb-build/DEBIAN/postinst
#!/usr/bin/env sh

mkdir -p /var/opt/hyperdeck-monitor/

systemctl daemon-reload
systemctl enable hyperdeck-monitor.service
systemctl restart hyperdeck-monitor.service
